# debian/rdiffweb.service
[Unit]
Description=Rdiffweb server
Documentation=man:rdiffweb(1)

[Service]
ExecStart=/usr/bin/rdiffweb

[Install]
WantedBy=multi-user.target
