# debian/rdiffweb.service
[Unit]
Description=Rdiffweb server

[Service]
ExecStart=/usr/bin/rdiffweb

[Install]
WantedBy=multi-user.target
